# patches/rootfs/etc/skel/.config/autostart/firstboot_user.desktop
[Desktop Entry]
Name=firstboot_user
Exec=firstboot_user
Terminal=false
Type=Application
Categories=Application;Utility
