# patches/rootfs/etc/profile.d/99bash.sh
export HISTCONTROL=ignoreboth
alias su-='su -'
alias rba='rpmbuild -ba'
alias rbs='rpmbuild -bs'
