# patches/rootfs/etc/xdg/autostart/xset_disable.desktop
[Desktop Entry]
Name=xset_disable
Exec=xset_disable
Terminal=false
Type=Application
Categories=Application;Utility
